# packaging/systemd/freeyourdisk.timer
[Unit]
Description=Weekly FreeYourDisk cache cleanup
Documentation=https://github.com/QrCommunication/FreeYourDisk

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
